# profile.d/git.sh
GIT_PROMPT_THEME=Default_Arch
source /usr/lib/bash-git-prompt/gitprompt.sh
eval "$(hub alias -s)"
