# profile.d/go.sh
export GOPATH=/var/shared
export PATH=$PATH:$GOPATH/bin
