# profile.d/ruby.sh
# replace `bundle exec`
export RUBYGEMS_GEMDEPS=-
