# profile.d/stardict.sh
export STARDICT_DIR="${HOME}/.config/nvim/plugged/vim-stardict"
source "${STARDICT_DIR}"/bindings/bash/stardict.sh

alias def="stardict"
alias defv="vstardict"

export STARDICT_PYTHON_PATH=/usr/bin/python2
