# profile.d/vim.sh
export XDG_CONFIG_HOME="$HOME/.config"
export EDITOR=nvim
alias vim=nvim
alias vi=nvim
